# files/SMPPSim.catcher.ini
[program:SMPPSim.catcher]
process_name = master
command=dotnet /opt/SmppSimCatcher/SmppSimCatcher.dll --smppsimcapturefile /opt/SMPPSim/sme_decoded.capture --urls "http://0.0.0.0:2875"
startsecs=0
autorestart=true