# files/SMPPSim.ini
[program:SMPPSim]
process_name = master
command=/etc/init.d/SMPPSim start
startsecs=0
autorestart=false