# files/SMPPSim.ini
[program:SMPPSim]
process_name = master
command=/opt/bin/SMPPSim.sh /opt/SMPPSim/conf/SMPPSim.props
startsecs=0
autorestart=true